# debian/tirex-nagios-plugin.postinst
#! /bin/sh

set -e

if [ -x /etc/init.d/nagios-nrpe-server ]
then
    if which invoke-rc.d >/dev/null 2>&1; then
        invoke-rc.d nagios-nrpe-server restart
    else
        /etc/init.d/nagios-nrpe-server restart
    fi
fi

#DEBHELPER#
