# debian/tirex-nagios-plugin.prerm
#! /bin/sh

set -e

dir=/etc/nagios/nrpe.d 
if [ -d $dir ]; then
    rm -f $dir/tirex-health.cfg
    rm -f $dir/tirex-queue-size.cfg
fi

if [ -x /etc/init.d/nagios-nrpe-server ]
then
    if which invoke-rc.d >/dev/null 2>&1; then
        invoke-rc.d nagios-nrpe-server restart
    else
        /etc/init.d/nagios-nrpe-server restart
    fi
fi

#DEBHELPER#
